# references/openclaw-skill-scan.path
[Unit]
Description=OpenClaw: watch ~/.openclaw/skills and trigger a scan

[Path]
PathModified=%h/.openclaw/skills
PathChanged=%h/.openclaw/skills

[Install]
WantedBy=default.target
